# res/conf/supervisor.ini
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
pidfile=/var/run/supervisord.pid 

[program:php-fpm]
command=sh /startPhpFpm.sh
autostart=true
autorestart=true

[program:nginx]
command=sh /startNginx.sh
autostart=true
autorestart=true

[program:v2ray]
command=sh /startV2Ray.sh
autostart=true
autorestart=true
