# res/conf/supervisord/supervisor.ini
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/var/tmp/supervisor/supervisord.pid 

[program:nginx]
command=sh /res/scripts/startNginx.sh
autostart=true
autorestart=true

[program:v2ray]
command=sh /res/scripts/startV2Ray.sh
autostart=true
autorestart=true
